# roles/personal_workstations/files/dropbox.service
[Unit]
Description=Dropbox Daemon

[Service]
ExecStart=flatpak run com.dropbox.Client

[Install]
WantedBy=default.target
