# roles/personal_workstations/files/syncthing.service
[Unit]
Description=Syncthing Daemon

[Service]
ExecStart=/usr/bin/syncthing serve --no-browser --logfile=default

[Install]
WantedBy=default.target
